# api/deploy/supervisor.conf
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; DR-PHARMA — Supervisor : Queue Workers
; Chemin sur le VPS : /etc/supervisor/conf.d/drpharma.conf
; Recharger : supervisorctl reread && supervisorctl update
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;

; ---------------------------------------------------------
; Worker 1 : queue par défaut + paiements (critique)
; On lance 2 processus parallèles pour absorber les pics
; ---------------------------------------------------------
[program:drpharma-worker-payments]
process_name     = %(program_name)s_%(process_num)02d
command          = php /var/www/drpharma/artisan queue:work database
                   --queue=payments,default
                   --sleep=3
                   --tries=3
                   --backoff=10
                   --max-time=3600
                   --memory=256
user             = drpharma
numprocs         = 2
autostart        = true
autorestart      = true
stopasgroup      = true
killasgroup      = true
startsecs        = 5
stopwaitsecs     = 60
redirect_stderr  = true
stdout_logfile   = /var/www/drpharma/storage/logs/worker-payments.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 3

; ---------------------------------------------------------
; Worker 2 : notifications (SMS, WhatsApp, FCM, email)
; ---------------------------------------------------------
[program:drpharma-worker-notifications]
command          = php /var/www/drpharma/artisan queue:work database
                   --queue=notifications
                   --sleep=5
                   --tries=3
                   --backoff=30
                   --max-time=3600
                   --memory=128
user             = drpharma
numprocs         = 1
autostart        = true
autorestart      = true
stopasgroup      = true
killasgroup      = true
startsecs        = 5
stopwaitsecs     = 60
redirect_stderr  = true
stdout_logfile   = /var/www/drpharma/storage/logs/worker-notifications.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 3

; ---------------------------------------------------------
; Groupe pour les redémarrer ensemble au déploiement
; ---------------------------------------------------------
[group:drpharma-workers]
programs = drpharma-worker-payments,drpharma-worker-notifications
