# deploy/drpharma-scheduler.service
# /etc/systemd/system/drpharma-scheduler.service
[Unit]
Description=DR-PHARMA Laravel Scheduler
After=network.target mysql.service

[Service]
Type=oneshot
User=drpharma
WorkingDirectory=/var/www/drpharma
ExecStart=/usr/bin/php8.3 artisan schedule:run
StandardOutput=append:/var/www/drpharma/storage/logs/scheduler.log
StandardError=append:/var/www/drpharma/storage/logs/scheduler.log
